# flights.cron
# China Airline
00 05 * * * /data/tickets/Flights/ChinaAirline/chinaairline.sh /data/tickets/Flights/ChinaAirline/asis.cfg
00 07 * * * /data/tickets/Flights/ChinaAirline/chinaairline.sh /data/tickets/Flights/ChinaAirline/europe.cfg
00 09 * * * /data/tickets/Flights/ChinaAirline/chinaairline.sh /data/tickets/Flights/ChinaAirline/america.cfg

# 酷航
00 03 * * * /data/tickets/Flights/Flyscoot/flyscoot.sh /data/tickets/Flights/Flyscoot/origin.cfg /data/tickets/Flights/Flyscoot/destination.cfg

# 捷星航空
00 01 * * * /data/tickets/Flights/JetStar/jetstar.sh /data/tickets/Flights/JetStar/destination.cfg

# 中國南方航空
#00 02 * * * /data/tickets/Flights/Csair/csair.sh /data/tickets/Flights/Csair/china.cfg /data/tickets/Flights/Csair/foreign.cfg

# 中國東方航空
#10 00 * * * /data/tickets/Flights/ChinaEastern/chinaeastern.sh /data/tickets/Flights/ChinaEastern/fromCity.cfg /data/tickets/Flights/ChinaEastern/asis.cfg
#10 06 * * * /data/tickets/Flights/ChinaEastern/chinaeastern.sh /data/tickets/Flights/ChinaEastern/fromCity.cfg /data/tickets/Flights/ChinaEastern/america.cfg
#10 12 * * * /data/tickets/Flights/ChinaEastern/chinaeastern.sh /data/tickets/Flights/ChinaEastern/fromCity.cfg /data/tickets/Flights/ChinaEastern/europe.cfg
#10 18 * * * /data/tickets/Flights/ChinaEastern/chinaeastern.sh /data/tickets/Flights/ChinaEastern/fromCity.cfg /data/tickets/Flights/ChinaEastern/oceanian.cfg

# 樂桃航空
#30 21 * * * /data/tickets/Flights/Peach/peach.sh /data/tickets/Flights/Peach/asis.cfg

# 香草航空
#00 22 * * * /data/tickets/Flights/VanillaAir/vanillaair.sh /data/tickets/Flights/VanillaAir/origin.cfg /data/tickets/Flights/VanillaAir/destination.cfg

# 長榮航空
00 06 * * * /data/tickets/Flights/EvaAirline/evaairline.sh /data/tickets/Flights/EvaAirline/asis.cfg
00 08 * * * /data/tickets/Flights/EvaAirline/evaairline.sh /data/tickets/Flights/EvaAirline/america.cfg
00 13 * * * /data/tickets/Flights/EvaAirline/evaairline.sh /data/tickets/Flights/EvaAirline/oceanian.cfg
00 15 * * * /data/tickets/Flights/EvaAirline/evaairline.sh /data/tickets/Flights/EvaAirline/europe.cfg
00 21 * * * /data/tickets/Flights/EvaAirline/evaairline.sh /data/tickets/Flights/EvaAirline/china.cfg

# 虎航
00 04 * * * /data/tickets/Flights/TigerAir/tigerair.sh /data/tickets/Flights/TigerAir/fromCity.cfg /data/tickets/Flights/TigerAir/toCity.cfg
